# scripts/com.machinecraft.ira-vitals.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.machinecraft.ira-vitals</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>/Users/rushabhdoshi/Desktop/Ira/scripts/ira_vital_signs.py</string>
        <string>--telegram</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/rushabhdoshi/Desktop/Ira/logs/launchd_vitals_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/rushabhdoshi/Desktop/Ira/logs/launchd_vitals_stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/rushabhdoshi/Desktop/Ira</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
